# docker/services/cron/run
#!/usr/bin/with-contenv sh

touch /etc/crontab /etc/cron.*/*

exec cron -f