# .bztrs/build/tasks/rootfs/etc/services.d/tasks/run
#!/usr/bin/with-contenv sh

exec /www/console scheduled-task:run --memory-limit=256M --time-limit=3600 --no-debug
