# docker/services/queue/run
#!/usr/bin/with-contenv sh

exec /www/bin/console scheduled-task:run --memory-limit=128M
