# .bztrs/build/queue/rootfs/etc/services.d/queue/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid www-data /www/bin/console messenger:consume default --memory-limit=128M --time-limit=60 --no-debug
