# .bztrs/build/queue/rootfs/etc/services.d/queue/run
#!/usr/bin/with-contenv sh

exec /www/console messenger:consume default --memory-limit=256M --time-limit=3600 --no-debug
